# dog_live_kd.service
[Unit]
Description=Dog Monitoring Live Stream Service
After=network.target

[Service]
Type=simple
ExecStart=/data/door-monitoring/start_dog_live_kd.sh
WorkingDirectory=/data/door-monitoring
Restart=always
RestartSec=5
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target